# systemd/lora-osmnotes.service
[Unit]
Description=Meshtastic → OSM Notes Gateway
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# SECURITY: User must be set to a non-root user
# The install script (scripts/install_pi.sh) will automatically detect and set the correct user
# If installing manually, replace 'nobody' with your username (e.g., User=pi or User=ubuntu)
# WARNING: Do not run as root! Change this to a non-privileged user.
User=nobody
Group=dialout
WorkingDirectory=/var/lib/lora-osmnotes
Environment="PATH=/opt/lora-osmnotes/bin"
Environment="DATA_DIR=/var/lib/lora-osmnotes"
ExecStart=/opt/lora-osmnotes/bin/python -m gateway.main
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=lora-osmnotes

# Security
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
